# templates/etc/init/ghost.conf.j2
# This file was generated by Ansible for {{ansible_fqdn}}
# Do NOT modify this file by hand!

start on startup

script
    cd {{ ghost_install_dir }}
    npm start {% if ghost_production %}--production{% endif %}

end script
